# templates/etc/cron.d/borgmatic.j2
# This file is generated by Ansible
#  DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
0 3 * * * {{ borgmatic_user_name }} {{ borgmatic_pip_virtualenv | default('/usr/local') }}/bin/borgmatic
